# files/etc/crontabs/root
* * * * * /root/startexec.sh & > /dev/null
*/5 * * * * /root/wifi/startscan.sh & > /dev/null
